# deploy/uq-proxy.service
[Unit]
Description=ProxmoxUserQuota transparent quota proxy (uq-proxy)
Documentation=https://github.com/WilliamLi0623/ProxmoxUserQuota-Proxy
After=network-online.target
Wants=network-online.target

[Service]
# Adjust paths/flags for your environment (see README).
ExecStart=/usr/local/bin/uq-proxy \
  -listen :8006 \
  -upstream https://CLOUD-ENTRY-NODE:8006 \
  -tls-cert /etc/uq-proxy/tls.crt \
  -tls-key /etc/uq-proxy/tls.key \
  -upstream-ca /etc/uq-proxy/pve-ca.pem
DynamicUser=yes
# Under DynamicUser the service runs as a transient user, so the TLS key
# must be reachable by group. Provision tls.key group-readable by www-data
# (e.g. acme.sh deploy, or reusing PVE's pve-ssl.key: chgrp www-data + chmod
# 640) and join that group here. Without this the proxy cannot read its key.
SupplementaryGroups=www-data
Restart=always
RestartSec=2
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/etc/uq-proxy

[Install]
WantedBy=multi-user.target
